# deploy/systemd/packaging-system.service
[Unit]
Description=Packaging System Node Service
After=network.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=PORT=8080
Environment=FORCE_HTTPS=0
ExecStart=__NODE_BIN__ __APP_DIR__/src/server.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
